# deploy/tbottabak.service
[Unit]
Description=TobaccoBot Telegram bot (oshisha.cc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tbottabak
Group=tbottabak
WorkingDirectory=/opt/tbottabak
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/tbottabak/.env
ExecStart=/opt/tbottabak/.venv/bin/python /opt/tbottabak/scripts/run_bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
